# banner-bot.service
[Unit]
Description=Banner Bot Discord Bot
After=network.target

[Service]
Type=simple
User=lucat
WorkingDirectory=/home/lucat/banner-bot
ExecStart=/usr/bin/node bot.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target 